# systemd/rmm-agent.service
# RMM Agent systemd Service Unit
#
# This unit file configures the RMM agent to run as a systemd Type=notify service.
# The agent notifies systemd when it's ready and sends periodic watchdog pings.
#
# Installation: cp rmm-agent.service /etc/systemd/system/
#               systemctl daemon-reload
#               systemctl enable rmm-agent.service
#
# Configuration: /etc/rmm-agent/config.yaml
# Logs: journalctl -u rmm-agent -f
#
# The agent runs as root to enable full system management capabilities:
# - Process management (kill any process)
# - Package installation/updates
# - System configuration changes
# - Service management

[Unit]
Description=RMM Agent - Remote Monitoring and Management
Documentation=https://github.com/doughall/linuxrmm
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStart=/opt/rmm-agent/bin/rmm-agent
Restart=on-failure
RestartSec=10s
StartLimitBurst=4
StartLimitInterval=5m

# Runtime directories (systemd creates these automatically)
RuntimeDirectory=rmm-agent
ConfigurationDirectory=rmm-agent
StateDirectory=rmm-agent
LogsDirectory=rmm-agent

# Watchdog - systemd kills service if no ping in 90 seconds
# Agent must send WATCHDOG=1 every interval/2 (~45 seconds)
WatchdogSec=90s

# Notify access - allow main process to send sd_notify
NotifyAccess=main

# Resource limits
LimitNOFILE=65535

# Environment
Environment=HOME=/var/lib/rmm-agent

[Install]
WantedBy=multi-user.target
